# docker/dev/supervisord/sshd.conf
[program:sshd]
directory=/home/dev
command=/usr/sbin/sshd -D
